# packages/duckbotos-branding/debian/source/etc/profile.d/duckbotos-branding.sh
#!/bin/sh
# Source the duckbotos branding if present
[ -r /etc/duckbotos/branding ] && . /etc/duckbotos/branding
