# packages/duckbotos-openclaw/debian/postinst
#!/bin/bash
# postinst for duckbotos-openclaw
# Installs OpenClaw gateway from npm + wires duckbot-memory brain plugin

set -e

REAL_HOME=$(getent passwd | awk -F: '$3 >= 1000 {print $6; exit}')
OPENCLAW_DIR="/opt/openclaw"
# OpenClaw state/config lives at ~/.openclaw/ per src/config/paths.ts
OPENCLAW_STATE_DIR="$REAL_HOME/.openclaw"
OPENCLAW_CONFIG_PATH="$OPENCLAW_STATE_DIR/openclaw.json"
SYSTEMD_USER_DIR="$REAL_HOME/.config/systemd/user"
DUCKBOTOS_MODES_DIR="/etc/duckbotos"
BRAIN_EXT_SRC="/opt/duckbotos/brain/extensions/duckbot-memory"
OPENCLAW_EXT_DIR="$OPENCLAW_STATE_DIR/extensions"

install_openclaw() {
    local home="${1:-$REAL_HOME}"

    echo "[duckbotos-openclaw] Installing OpenClaw..."

    # Create directories
    mkdir -p "$OPENCLAW_DIR"
    mkdir -p "$OPENCLAW_STATE_DIR"
    mkdir -p "$SYSTEMD_USER_DIR"
    mkdir -p "$OPENCLAW_EXT_DIR"

    # Clone OpenClaw if not present
    if [ ! -d "$OPENCLAW_DIR/.git" ]; then
        echo "[duckbotos-openclaw] Cloning OpenClaw from openclaw/openclaw..."
        rm -rf "$OPENCLAW_DIR"
        git clone --depth=1 --branch main \
            https://github.com/openclaw/openclaw \
            "$OPENCLAW_DIR"
    fi

    # Install dependencies and build
    echo "[duckbotos-openclaw] Installing npm dependencies..."
    cd "$OPENCLAW_DIR"
    npm install --omit=dev 2>&1 | tail -5 || true
    npm run build 2>&1 | tail -5 || true

    # Create openclaw CLI wrapper
    cat > /usr/local/bin/openclaw << 'OPENCLAW_WRAPPER'
#!/bin/bash
# OpenClaw CLI wrapper for DuckBotOS
export OPENCLAW_HOME="${OPENCLAW_HOME:-${REAL_HOME}/.openclaw}"
export OPENCLAW_DATA="${OPENCLAW_DATA:-${REAL_HOME}/.openclaw/workspace}"
export NODE_ENV=production
cd /opt/openclaw
exec node /opt/openclaw/dist/index.js "$@"
OPENCLAW_WRAPPER
    chmod +x /usr/local/bin/openclaw

    # Install systemd user service
    cat > "$SYSTEMD_USER_DIR/openclaw-gateway.service" << SERVICE_EOF
[Unit]
Description=OpenClaw Gateway
Documentation=https://github.com/openclaw/openclaw
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=${OPENCLAW_STATE_DIR}
Environment=OPENCLAW_HOME=${OPENCLAW_DIR}
Environment=OPENCLAW_STATE_DIR=${OPENCLAW_STATE_DIR}
Environment=NODE_ENV=production
Environment=PORT=18789
ExecStart=/usr/local/bin/openclaw gateway start --port 18789 --workspace ${OPENCLAW_STATE_DIR}
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=${OPENCLAW_STATE_DIR} ${OPENCLAW_DIR}
PrivateTmp=true

[Install]
WantedBy=default.target
SERVICE_EOF

    # Set permissions
    chown -R root:root /opt/openclaw
    chown -R root:root /usr/local/bin/openclaw
    chmod -R 755 /opt/openclaw

    # Enable lingering so service starts at boot without login
    if command -v loginctl >/dev/null 2>&1; then
        local user=$(who | grep -v '^root ' | head -1 | awk '{print $1}')
        loginctl enable-linger "$user" 2>/dev/null || true
    fi

    # Enable and start the service
    if command -v systemctl >/dev/null 2>&1; then
        systemctl --user daemon-reload 2>/dev/null || true
        systemctl --user enable openclaw-gateway.service 2>/dev/null || true
        systemctl --user start openclaw-gateway.service 2>/dev/null || true
    fi

    echo "[duckbotos-openclaw] OpenClaw installed."
}

register_brain_plugin() {
    # Wire the duckbot-memory brain plugin (installed by duckbotos-brain package)
    echo "[duckbotos-openclaw] Registering duckbot-memory brain plugin..."

    # Copy the OpenClaw native plugin shim into OpenClaw's extensions dir
    if [ -d "$BRAIN_EXT_SRC" ]; then
        cp -r "$BRAIN_EXT_SRC" "$OPENCLAW_EXT_DIR/duckbot-memory"
        echo "[duckbotos-openclaw] Brain plugin copied to $OPENCLAW_EXT_DIR/duckbot-memory"
    else
        echo "[duckbotos-openclaw] Warning: duckbotos-brain not yet installed. Skipping brain plugin."
        return 0
    fi

    # Register plugin in OpenClaw config via python (safe JSON patching)
    # NOTE: OpenClaw config lives at $HOME/.openclaw/openclaw.json per
    # src/config/paths.ts:resolveCanonicalConfigPath returns ~/.openclaw/openclaw.json.
    # NOT /var/lib/openclaw/workspace/openclaw.json (that was wrong).
    OPENCLAW_CFG="$REAL_HOME/.openclaw/openclaw.json"
    mkdir -p "$REAL_HOME/.openclaw"
    if [ -n "$REAL_USER" ]; then
        chown -R "$REAL_USER:$REAL_USER" "$REAL_HOME/.openclaw" 2>/dev/null || true
    fi
    if [ ! -f "$OPENCLAW_CFG" ]; then
        echo '{}' > "$OPENCLAW_CFG"
        [ -n "$REAL_USER" ] && chown "$REAL_USER:$REAL_USER" "$OPENCLAW_CFG" 2>/dev/null || true
    fi

    python3 - << PYEOF
import json, os, pwd
real_user = '$REAL_USER'
real_home = '$REAL_HOME'
cfg_path = f'{real_home}/.openclaw/openclaw.json'
with open(cfg_path) as f:
    cfg = json.load(f)

if 'plugins' not in cfg:
    cfg['plugins'] = {}
if 'entries' not in cfg['plugins']:
    cfg['plugins']['entries'] = {}

cfg['plugins']['entries']['duckbot-memory'] = {
    'enabled': True,
    'config': {
        'repoPath': '/opt/duckbotos/brain',
        'autoWakeUp': True,
        'autoSync': True,
        'defaultK': 8
    }
}

with open(cfg_path, 'w') as f:
    json.dump(cfg, f, indent=2)
try:
    os.chown(cfg_path, pwd.getpwnam(real_user).pw_uid, -1)
except Exception:
    pass
print(f'[duckbotos-openclaw] duckbot-memory registered at {cfg_path}')
PYEOF
}

with open(cfg_path, 'w') as f:
    json.dump(cfg, f, indent=2)
print('[duckbotos-openclaw] duckbot-memory registered in openclaw.json')
PYEOF
}

configure_kiosk_mode() {
    # Configure the kiosk to load OpenClaw by default
    mkdir -p "$DUCKBOTOS_MODES_DIR"
    echo "http://localhost:18789/plugins/openclawos" > "$DUCKBOTOS_MODES_DIR/kiosk-mode"
    echo "[duckbotos-openclaw] Kiosk mode set to OpenClaw workspace."
}

case "$1" in
    configure)
        install_openclaw
        register_brain_plugin
        configure_kiosk_mode

        echo ""
        echo "✅ OpenClaw installed for DuckBotOS"
        echo "   Gateway:  http://localhost:18789"
        echo "   Brain:    duckbot-memory plugin (67 tools + session hooks)"
        echo "   Config:   $OPENCLAW_DIR"
        echo ""
        echo "   To restart: systemctl --user restart openclaw-gateway"
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        ;;
esac