# packages/duckbotos-session-picker/session-picker.service
[Unit]
Description=DuckBotOS Session Picker
Documentation=https://github.com/Franzferdinan51/DuckBotOS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/duckbotos-session-picker
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/etc/duckbotos

[Install]
WantedBy=multi-user.target