# packages/duckbotos-lm-studio/debian/postinst
#!/bin/bash
# postinst for duckbotos-lm-studio — configure LM Studio API server

set -e

LM_STUDIO_SERVICE="/etc/systemd/system/lm-studio-api.service"
LM_STUDIO_CONFIG_DIR="/etc/lm-studio"
LM_STUDIO_CONFIG="${LM_STUDIO_CONFIG_DIR}/config.yaml"

case "$1" in
    configure)
        # Create config directory
        mkdir -p "$LM_STUDIO_CONFIG_DIR"

        # Write default config (localhost-only API, no auth for local users)
        cat > "$LM_STUDIO_CONFIG" << 'CONFIG'
# LM Studio configuration for DuckBotOS
# Generated by duckbotos-lm-studio package

server:
  host: "127.0.0.1"    # localhost only — security
  port: 1234           # LM Studio API port
  # auth_token: ""    # Set via environment LM_STUDIO_AUTH_TOKEN

paths:
  models: "/var/lib/lm-studio/models"   # Shared model directory
  downloads: "/var/lib/lm-studio/downloads"

gpu:
  enabled: true
  # auto: detect GPU automatically
  nvidia: true    # CUDA (NVIDIA GPUs)
  amd: false      # ROCm (AMD GPUs — set true if AMD GPU present)
  vulkan: false   # Vulkan fallback

model_defaults:
  context_length: 4096
  gpu_offload: true
  threads: 4
CONFIG

        # Set permissions
        chmod 644 "$LM_STUDIO_CONFIG"

        # Enable and start LM Studio API service
        if [ -f "$LM_STUDIO_SERVICE" ] && command -v systemctl >/dev/null 2>&1; then
            systemctl enable lm-studio-api.service 2>/dev/null || true
            systemctl start lm-studio-api.service 2>/dev/null || true
        fi
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
    *)
        ;;
esac